# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:vnc]
command=/bin/bash -c "su - vncuser -c 'vncserver :1 -geometry 1920x1080 -depth 24 -fg'"
autostart=true
autorestart=true
stderr_logfile=/var/log/vnc.err.log
stdout_logfile=/var/log/vnc.out.log
priority=100
startsecs=10
startretries=3

[program:websockify]
command=websockify --web /usr/share/novnc 6080 localhost:5901
autostart=true
autorestart=true
stderr_logfile=/var/log/websockify.err.log
stdout_logfile=/var/log/websockify.out.log
priority=200
startsecs=5
startretries=3

# [program:nginx]
# command=nginx -g "daemon off;"
# autostart=true
# autorestart=true
# stderr_logfile=/var/log/nginx.err.log
# stdout_logfile=/var/log/nginx.out.log
# priority=300
# startsecs=5
# startretries=3